# BaseImages/nginx-php/local/src/resources/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -e -f /home/sindria/custom_ssh/sshd_config
numprocs=1
autostart=true
autorestart=true
user=sindria