# BaseImages/nginx-php/main/src/resources/supervisor/conf.d/php-fpm-log-access.conf
[program:php-fpm-access-log]
command=/usr/bin/tail -f /var/log/php/access.log
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true
startsecs=3