# Core/Base/nginx-php/main/src/resources/supervisor/conf.d/php-fpm-log-error.conf
[program:php-fpm-error-log]
command=/usr/bin/tail -f /var/log/php/error.log
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true
startsecs=3